# root/etc/cron.hourly/dance-hourly
#!/bin/sh
## This script is run every hour.
dance-fetch-students
dance-fetch-calendar
dance-update-todays-classes  # Just in case...
dance-upload-to-spreadsheet
